# drafts/tests/supervisord.conf
[supervisord]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:dockerd]
command = bash -c "/usr/bin/dockerd; kill -s SIGINT `cat supervisord.pid`"
priority=1
autostart=true
autorestart=false
startsecs=0
stopasgroup = true

[program:nix-daemon]
command = bash -c "/root/.nix-profile/bin/nix-daemon; kill -s SIGINT `cat supervisord.pid`"
priority=1
autostart=true
autorestart=false
startsecs=0
stopasgroup = true

[program:nix-serve]
command = bash -c "/root/.nix-profile/bin/nix-serve --listen 0.0.0.0:8000; kill -s SIGINT `cat supervisord.pid`"
priority=1
autostart=true
autorestart=false
startsecs=0
stopasgroup = true
